# systemd/user/device-sounds.service
[Unit]
Description=Device connect/disconnect sounds

[Service]
ExecStart=%h/.config/scripts/device_sounds.sh
Restart=always
